# src/web_ui/resources/config/crontab
# Bypass Keenetic - Cron schedule
# Update ipset every 6 hours
0 */6 * * * /opt/bin/unblock_update.sh
